# files/sr2/etc/cron.d/sr2-monthly
30 1 1 * * root /usr/local/sr2/.venv/bin/python /usr/local/sr2/app/python/batch/monthly_record_count.py -c /usr/local/sr2/app/stat.secret.json >> /var/log/sr2-cron-monthly.log
